# apps/io/keyd/keyd.service
[Unit]
Description=key remapping daemon

[Service]
Type=simple
ExecStart=keyd
